# debian/cinelerra.postrm
#! /bin/sh
# Ubuntu Configure script for cinelerra by Rampino Paolo under GNU Public License.

set -e

#DEBHELPER#
# Define
IC=cinestart
SCRIPT=/etc/init.d/$IC
STARTCINELERRA=/usr/bin/Cinelerra

if [ "$1" = "remove" ] ; then
rm $STARTCINELERRA
/usr/sbin/update-rc.d -f $IC remove
rm $SCRIPT

fi

exit 0


